# rootfs/etc/cron.d/counters
00 00 * * * root echo 0 > /tmp/counters/_parent
00 00 * * * root echo 0 > /tmp/counters/cron
00 00 * * * root echo 0 > /tmp/counters/apache
