# rootfs/etc/cron.d/sendy
*/01 * * * * root /usr/local/bin/php /var/www/html/import-csv.php > /dev/null 2>&1
*/01 * * * * root /usr/local/bin/php /var/www/html/autoresponders.php > /dev/null 2>&1
*/05 * * * * root /usr/local/bin/php /var/www/html/scheduled.php > /dev/null 2>&1
*/15 * * * * root /usr/local/bin/php /var/www/html/update-segments.php > /dev/null 2>&1
